# scripts/com.mindquest.n8n-backup.plist
no such file or directory: scripts/run-n8n-backup.mjs
aldosantos@unknownfe5a1daca9e2 MindQuest % <?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.mindquest.n8n-backup</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>/Users/aldosantos/Documents/GitHub/MindQuest/scripts/run-n8n-backup.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
      <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/aldosantos/Library/Logs/mindquest-n8n-backup.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/aldosantos/Library/Logs/mindquest-n8n-backup.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
